# workbench-frontend.service
[Unit]
Description=Workbench Frontend
After=network.target

[Service]
Type=simple
User=workbench
Group=workbench
WorkingDirectory=/opt/workbench/frontend
ExecStart=/usr/bin/pnpm preview --port 4510 --host 0.0.0.0
Restart=always
RestartSec=10
Environment="NODE_ENV=production"
EnvironmentFile=-/opt/workbench/.env

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=append:/var/log/workbench/frontend.log
StandardError=append:/var/log/workbench/frontend-error.log

[Install]
WantedBy=multi-user.target