# config/workbench-backend.service
[Unit]
Description=Workbench Backend API
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=workbench
Group=workbench
WorkingDirectory=/opt/workbench/backend
ExecStart=/opt/workbench/backend/workbench-server
Restart=always
RestartSec=10
EnvironmentFile=/opt/workbench/.env

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/mnt/nas /var/log/workbench

# Logging
StandardOutput=append:/var/log/workbench/backend.log
StandardError=append:/var/log/workbench/backend-error.log

[Install]
WantedBy=multi-user.target